# ops/com.agentchappie.watchdog.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.agentchappie.watchdog</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/env</string>
    <string>python3</string>
    <string>/Users/chappie/Projects/Agent.Chappie/scripts/watchdog_agent.py</string>
    <string>--status-dir</string>
    <string>/Users/chappie/Projects/Agent.Chappie/runtime_status</string>
    <string>--stale-seconds</string>
    <string>60</string>
    <string>--label</string>
    <string>com.agentchappie.runtime</string>
    <string>--max-restarts</string>
    <string>3</string>
    <string>--window-seconds</string>
    <string>300</string>
    <string>--check-queue-consumer</string>
    <string>--remediate-duplicate-consumers</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>StartInterval</key>
  <integer>30</integer>
  <key>KeepAlive</key>
  <false/>
  <key>StandardOutPath</key>
  <string>/Users/chappie/Projects/Agent.Chappie/runtime_status/watchdog_stdout.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/chappie/Projects/Agent.Chappie/runtime_status/watchdog_stderr.log</string>
  <key>WorkingDirectory</key>
  <string>/Users/chappie/Projects/Agent.Chappie</string>
  <key>ThrottleInterval</key>
  <integer>30</integer>
</dict>
</plist>
